# deploy/systemd/thulium-bridge.service
[Unit]
Description=Thulium Bridge Worker (loop mode)
After=network.target mysql.service

[Service]
Type=simple
WorkingDirectory=/opt/thulium-bridge
ExecStart=/usr/bin/php /opt/thulium-bridge/bin/thulium-bridge run-loop
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
User=www-data
Group=www-data
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
