# deploy/cron/example.cron
* * * * * /usr/bin/php /opt/thulium-bridge/bin/thulium-bridge run-once >> /var/log/thulium-bridge-cron.log 2>&1
*/15 * * * * /usr/bin/php /opt/thulium-bridge/bin/thulium-bridge health >> /var/log/thulium-bridge-health.log 2>&1
0 3 * * * /usr/bin/php /opt/thulium-bridge/bin/thulium-bridge cleanup >> /var/log/thulium-bridge-cleanup.log 2>&1
